# rootfs/multistrap/rootfs/usr/share/upstart/sessions/ssh-agent.conf
description "SSH Agent"
author "Stéphane Graber <stgraber@ubuntu.com>"

start on starting xsession-init

exec /usr/lib/openssh/agent-launch start
post-stop exec /usr/lib/openssh/agent-launch stop
