# rootfs/multistrap/rootfs/etc/inittab
T0:2345:respawn:/sbin/getty -L ttyS0 115200 vt100
