# rootfs/multistrap/rootfs/var/lib/dpkg/info/dphys-swapfile.conffiles
/etc/dphys-swapfile
/etc/init.d/dphys-swapfile
