# rootfs/multistrap/rootfs/var/lib/dpkg/info/dphys-swapfile.list
/.
/etc
/etc/dphys-swapfile
/etc/init.d
/etc/init.d/dphys-swapfile
/sbin
/sbin/dphys-swapfile
/usr
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/dphys-swapfile
/usr/share/doc
/usr/share/doc/dphys-swapfile
/usr/share/doc/dphys-swapfile/FAQ
/usr/share/doc/dphys-swapfile/NEWS.Debian.gz
/usr/share/doc/dphys-swapfile/README
/usr/share/doc/dphys-swapfile/changelog.Debian.gz
/usr/share/doc/dphys-swapfile/changelog.gz
/usr/share/doc/dphys-swapfile/copyright
/usr/share/doc/dphys-swapfile/examples
/usr/share/doc/dphys-swapfile/examples/dphys-swapfile.example
/usr/share/doc/dphys-swapfile/examples/init.d.example
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/dphys-swapfile
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/dphys-swapfile.8.gz
