# rootfs/multistrap/rootfs/var/lib/dpkg/info/openssh-server.conffiles
/etc/default/ssh
/etc/init.d/ssh
/etc/init/ssh.conf
/etc/network/if-up.d/openssh-server
/etc/pam.d/sshd
/etc/ufw/applications.d/openssh-server
